# files/etc/systemd/system/var-lib-pacman.mount
# Added by https://github.com/logan2611/holo-overlay
# License: MIT License

[Unit]
Description=Overlay for /var/lib/pacman/local
After=etc.mount
After=opt.mount

[Mount]
Where=/var/lib/pacman/
What=overlay
Type=overlay
Options=lowerdir=/var/lib/pacman/,upperdir=/opt/holo-overlay/var/lib/pacman/upper,workdir=/opt/holo-overlay/var/lib/pacman/work,index=off,metacopy=off

[Install]
WantedBy=local-fs.target
